# tests/supervisor.conf
[supervisord]
logfile = %(here)s/var/supervisord.log
pidfile = %(here)s/var/supervisord.pid

[program:bro]
command = python %(here)s/../src/bro.py
directory = %(here)s

[program:pub]
command = sh -c "sleep 0.5 && python %(here)s/../src/pub.py foo %(here)s/data"
directory = %(here)s

[program:sub]
command = sh -c "sleep 1.0 && python %(here)s/../src/sub.py"
directory = %(here)s
