# tests/supervisor.conf
[supervisord]
logfile = %(here)s/supervisord.log
pidfile = %(here)s/supervisord.pid

[program:bro]
command = %(here)s/../venv/bin/python %(here)s/../src/bro.py

[program:pub]
command = bash -c "sleep 0.5 && %(here)s/../venv/bin/python %(here)s/../src/pub.py foo %(here)s/data"

[program:sub]
command = bash -c "sleep 1.0 && %(here)s/../venv/bin/python %(here)s/../src/sub.py"
