# tests/supervisor.conf
[supervisord]

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Unfortunately there is no easy way to avoid replicating the RPC socket path
# in a way that works for both Supervisor daemon and control client
# (nor of having it passed via the command line)
# -- at least not without using environment variables.
[unix_http_server]
file = %(here)s/var/supervisord.sock

[supervisorctl]
serverurl = unix://%(here)s/var/supervisord.sock

# The delays (sleep and startsecs) below are
# a hacky way to serialize the start of programs
# until the connections that they require can be made,
# but one that still keeps them independent from each other,
# so that we can check if each of them fails to start, crashes, etc.

[program:bro]
command = python %(here)s/../src/bro.py
directory = %(here)s
startsecs = 1

[program:pub]
command = sh -c "sleep 0.5 && python %(here)s/../src/pub.py foo %(here)s/data"
directory = %(here)s
startsecs = 2

[program:sub]
command = sh -c "sleep 1.0 && python %(here)s/../src/sub.py"
directory = %(here)s
startsecs = 3
